# snippets/admin/minecraft/files/crontab
@reboot echo It worked... > /root/worked.txt
